# deploy/roomworks/supervisord.conf
[supervisord]
nodaemon=true
pidfile=/tmp/supervisord.pid
user=root
loglevel=info

[program:signaling]
command=/usr/local/bin/roomworks-signaling
autostart=true
autorestart=true
startretries=10
priority=10
stopasgroup=true
killasgroup=true
umask=0077
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=10
priority=20
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
